# flight_player/udev/88-init_CAN.rules
SUBSYSTEM="usb", ACTION="add",ATRRS{idVendor}=="1d50",
ATTRS{idProduct}=="606f", RUN+="/usr/local/initialize_CAN.sh"
